# backend/gunicorn.service
[Unit]
Description=LiveStream Monitoring Backend Service (Gevent Optimized)
After=network.target
Wants=network-online.target                  # Ensure network is fully up
Documentation=https://live-stream-monitoring-vue3-flask.vercel.app/   # Link to service docs

[Service]
Type=simple
User=ec2-user
Group=ec2-user

WorkingDirectory=/home/ec2-user/LiveStream_Monitoring_Vue3_Flask/backend

# Load environment overrides (create /etc/sysconfig/livemon-backend)
Environment=PYTHONUNBUFFERED=1               # Unbuffered logs to stdout
Environment=FLASK_DEBUG=true
Environment=ENABLE_SSL=true                  # Enable SSL by default
Environment=CERT_DIR=/home/ec2-user/LiveStream_Monitoring_Vue3_Flask/backend
Environment=SSL_CERT_PATH=/home/ec2-user/LiveStream_Monitoring_Vue3_Flask/backend/fullchain.pem
Environment=SSL_KEY_PATH=/home/ec2-user/LiveStream_Monitoring_Vue3_Flask/backend/privkey.pem

# Gunicorn startup (no shell wrapper for cleaner exec)
ExecStart=/opt/pytorch/bin/gunicorn \
  --worker-class gevent \
  --workers 2 \
  --threads 4 \
  --timeout 180 \
  --bind 0.0.0.0:5000 \
  --certfile=${SSL_CERT_PATH} \
  --keyfile=${SSL_KEY_PATH} \
  --log-level info \
  main:app

# Restart on failure, with back-off
Restart=on-failure
RestartSec=5s
StartLimitBurst=3                            # Only try 3 restarts in a row
StartLimitInterval=60s                       # Within a 60s interval

# Resource limits & security hardening
LimitNOFILE=65536                            # Raise file descriptor limit
NoNewPrivileges=true                         # Prevent privilege escalation
PrivateTmp=true                              # Isolate /tmp directory
ProtectSystem=full                           # Mount /usr and /boot read-only
ProtectHome=read-only                        # Read-only access to home dirs

StandardOutput=journal                       # Send stdout to journal
StandardError=journal                        # Send stderr to journal
SyslogIdentifier=livemon-backend

[Install]
WantedBy=multi-user.target